# service.template/buffett-leader.socket
[Unit]
Description=Buffett leader socker
PartOf=buffett-leader.service

[Socket]
ListenStream=8899

[Install]
WantedBy=sockets.target